# supervisor/nodejs.conf
[program:nodejs]
user=agung
command=node /home/agung/nodejs.kelasdevops.xyz/server.js
autostart=true
autorestart=true
stderr_logfile=/var/log/nodejs.err.log
stdout_logfile=/var/log/nodejs.out.log